# docker/supervisor/my_resume.ini
[program:my_resume]
command=python3 /www/web/myresume/readdoc.py
directory=/www/web/myresume/
process_name=%(program_name)s
num_procs=1
priority=990
stopasgroup=true
autostart=true
autorestart=true
stdout_logfile_maxbytes=40MB
redirect_stderr=true
stopsignal=QUIT
stdout_logfile = /var/log/supervisor/myresume.log
